# system/homework-ovsdb-server.conf
start on stopped networking
stop on runlevel [016]
expect fork
exec /usr/local/sbin/ovsdb-server --remote=punix:/usr/local/var/run/openvswitch/db.sock --detach --monitor --pidfile
