# system/iptables-masquerade.conf
start on startup
stop on runlevel [06]
script
    . /home/homeuser/router.conf
    iptables -t nat -A POSTROUTING -o $GATEWAY -j MASQUERADE
end script
