# system/iptables-masquerade.conf
start on startup
stop on runlevel [016]
script
    iptables -t nat -A POSTROUTING -o $GATEWAY -j MASQUERADE
end script
