# system/ovs-vswitchd.conf
start on started ovsdb-server
stop on stopping ovsdb-server
expect fork
exec /usr/local/sbin/ovs-vswitchd --detach --monitor --pidfile
